# rootfs_overlay/etc/inittab
# /etc/inittab

::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mount -t sysfs sysfs /sys
::sysinit:/bin/mount -t tmpfs tmpfs /tmp
::sysinit:/bin/mount -o remount,rw /

::sysinit:/etc/init.d/rcS

console::respawn:/bin/su - autumnuser0 -c "/bin/sh"

::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/swapoff -a
::shutdown:/bin/umount -a -r
::restart:/sbin/init
